# automentohust.sh
#!/bin/sh /etc/rc.common

path=""
START=50
start()
{
  `sh $path` &
}
